# lib/health.service
[Unit]
Description=Port Health Check Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/health_check.sh